# docker/supervisor/programs/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm_err.log
